# systemd/ardupilot-swarm.service.in
[Unit]
Description=ArduPilot swarm
Wants=network-online.target mavlink-router.service
After=network-online.target mavlink-router.service

[Service]
Type=oneshot
User=@RUN_USER@
Group=@RUN_GROUP@
Environment=HOME=@RUN_HOME@
ExecStart=/usr/local/bin/start-ardupilot-swarm
ExecStop=/usr/local/bin/stop-ardupilot-swarm
RemainAfterExit=yes
TimeoutStartSec=60
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
